# deploy/systemd/smart-todolist-agent.service
[Unit]
Description=Smart TodoList - .NET Agent Backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/smart-todolist/current/agent
ExecStart=/var/www/smart-todolist/current/agent/ProverbsAgent
Restart=on-failure
RestartSec=5
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://127.0.0.1:8000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smart-todolist-agent

[Install]
WantedBy=multi-user.target
